# crontab
*/15 * * * * cd /home/zenbot/zenbot/ && ./zenbot.sh backfill gdax.ETH-EUR --days 15

0 0 * * * cd /home/zenbot/zenbot-insight/ && ./batch_strategies.sh
